# pyramid-exporter.timer
[Unit]
Description=Run Pyramid metrics collector every 15 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=15s
AccuracySec=1s

[Install]
WantedBy=timers.target
